# update-system.service
[Unit]
Description=Update the system
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/usr/sbin/update-system -s
